# scripts/launchd/com.prd-agent.monthly.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prd-agent.monthly</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_DIR}}/scripts/run-monthly.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>
    <key>StandardOutPath</key>
    <string>{{PROJECT_DIR}}/logs/monthly.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_DIR}}/logs/monthly-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.nvm/versions/node/v22.19.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>
</dict>
</plist>
